# home/.config/systemd/user/rclone-s3-documents.timer
[Unit]
Description=Run document bisync every 5 minutes
Wants=network-online.target
After=network-online.target

[Timer]
OnStartupSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=30

[Install]
WantedBy=default.target
